# systemd/export.service
[Unit]
Description=Export image

[Service]
Type=oneshot
WorkingDirectory=/opt/rtsp-exporter
ExecStart=/opt/rtsp-exporter/export.sh
