# systemd/pull.service
[Unit]
Description=Pull repository

[Service]
Type=oneshot
WorkingDirectory=/opt/rtsp-exporter
ExecStart=git pull
